# systemd/fr33d0m-neurovision.service
[Unit]
Description=Fr33d0m Neurovision Terminal Visualizer (daemon mode)
After=network.target

[Service]
Type=simple
User=%i
Environment=HERMES_HOME=%h/.hermes
Environment=PYTHONPATH=%h/.hermes/extensions/hermes-neurovision
ExecStart=%h/.hermes/hermes-agent/venv/bin/python -m hermes_neurovision.cli --daemon --logs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
