# systemd/fr33d0m-neurovision-web.service
[Unit]
Description=Fr33d0m Browser Neurovision
After=network.target

[Service]
Type=simple
User=%i
Environment=HERMES_HOME=%h/.hermes
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/ttyd -p 7682 -i lo -b /neurovision-proxy -t fontSize=14 -t cursorStyle=bar %h/.local/bin/fr33d0m-neurovision-shell
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
